# ops/autonomous-trading.service
[Unit]
Description=Rust Algorithmic Trading System - Autonomous Mode
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/Developer/RustAlgorithmTrading
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=%h/Developer/RustAlgorithmTrading/ops/scripts/autonomous_trading_system.sh --mode=continuous
Restart=always
RestartSec=60
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/Developer/RustAlgorithmTrading/logs
ReadWritePaths=%h/Developer/RustAlgorithmTrading/data

[Install]
WantedBy=multi-user.target
